# etc/systemd/system/camera.service
[Unit]
Description=Simple NMEA logger

[Service]
ExecStart=/home/pi/carcam/scripts/camera.sh
Restart=always

[Install]
WantedBy=multi-user.target

