# etc/systemd/system/nmea.service
[Unit]
Description=Simple NMEA logger

[Service]
ExecStart=/home/pi/carcam/scripts/nmea5.sh
Restart=always

[Install]
WantedBy=multi-user.target

